# contrib/systemd/panotify.service.in
[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=@bindir@/sysnotifier
Restart=always
RestartSec=10

[Unit]
After=graphical-session.target
Description=Key system events notifier
PartOf=graphical-session.target

